# applications/luci-app-aliddns/root/etc/init.d/aliddns
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 Lienol <lawlienol@gmail.com>

START=98

CONFIG=aliddns
CONFIG_PATH=/var/etc/$CONFIG
cronfile="/usr/ddns/cronfile"
tempfile="/usr/ddns/tempfile"


run_service() {
#	config_get enable $1 enable
#	[ "$enable" = "" ] && return 0
	config_get AccessKeyId $1 AccessKeyId
	[ "AccessKeyId" = "" ] && return 0
	config_get AccessKeySecret $1 AccessKeySecret
	[ "AccessKeySecret" = "" ] && return 0
	config_get Subdomain $1 Subdomain
	config_get Domain $1 Domain
	config_get Iptype $1 Iptype
	config_get TTL $1 TTL
	config_get DnsServer $1 DnsServer

	echo "*/5 * * * * /bin/bash /usr/ddns/ddns.sh $AccessKeyId $AccessKeySecret $Subdomain $Domain $Iptype $TTL $DnsServer" >> /usr/ddns/cronfile

}
conf_set(){
  keyword="ddns"
  rm -rf $cronfile
  rm -rf $tempfile
  `crontab -l > $tempfile`
  while read -r line
  do
    if [[ "$line" == "*$keyword*" ]]
    then
      pass
    else
      `echo "$line" >> $cronfile`
    fi
  done < $tempfile
}

start() {
#	enable=$(uci -q get $CONFIG.@global[0].enable)
#	if [ "$enable" = "0" ];then
#		stop_service
#	else
#		mkdir -p $CONFIG_PATH
#		rm -f /var/etc/$CONFIG.port
#		config_load $CONFIG
#		config_foreach run_service "config"
	conf_set
  config_load $CONFIG
  config_foreach run_service "alikey"
  crontab -r
  crontab $cronfile
#	fi
}

stop() {
  conf_set
  crontab -r
  crontab $cronfile
}

restart() {
	start
}
