# applications/luci-app-cpu/root/etc/uci-defaults/luci-app-cpu
#!/bin/sh


chmod +x /usr/libexec/rpcd/luci.cpu

exit 0
