# applications/luci-app-aliddns/root/etc/uci-defaults/40_luci-aliddns
#!/bin/sh

crontab /usr/ddns/cronfile

return 0
